# cnpj-api.service
[Unit]
Description=CNPJ API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/cnpj-api
Environment=PATH=/root/cnpj-api/.venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/root/cnpj-api/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 1
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
